# debian/postinst
systemctl enable com-ukui-desktop-software.service
