# .pc/change_the_config_of_dbus.patch/etc/com-ukui-desktop-software.service
[Unit]
Description=mdm service
After=dbus.service

[Service]
Type=dbus
BusName=com.kylin.mdm
ExecStart=/usr/bin/safectl
NotifyAccess=main

